# iptv-capture.service
[Unit]
Description=DTH IPTV Capture Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dth-capture
EnvironmentFile=/home/pi/dth-capture/.env
ExecStart=/home/pi/dth-capture/.venv/bin/python -m app.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
